# backend/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:api]
command=gunicorn app.main:app
    -k uvicorn.workers.UvicornWorker
    -w %(ENV_WEB_WORKERS)s
    --bind 0.0.0.0:8000
    --graceful-timeout 30
    --timeout 120
    --keep-alive 5
    --access-logfile -
    --error-logfile -
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:worker]
command=celery -A worker.tasks.celery_app worker
    --loglevel=info
    -c %(ENV_WORKER_CONCURRENCY)s
    --max-tasks-per-child=10
    --without-gossip
    --without-mingle
directory=/app
environment=PYTHONPATH=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
